# yota-auto.init
#!/bin/sh /etc/rc.common

START=66
BIN=yota-auto
DEFAULT=/etc/default/$BIN
RUN_D=/var/run
PID_F=$RUN_D/$BIN.pid

start() {
        [ -f $DEFAULT ] && . $DEFAULT
        mkdir -p $RUN_D
        $BIN
}

stop() {
        [ -f $PID_F ] && kill $(cat $PID_F)
        [ -f $PID_F ] && rm $PID_F
}
